# webapp/etc/systemd/isu-go.service.d/override.conf
# isu-go GCチューニング drop-in（本番反映先: /etc/systemd/system/isu-go.service.d/override.conf）
# 反映: sudo systemctl daemon-reload && sudo systemctl restart isu-go
[Service]
# 第14回: GCチューニング。app RSSは~21MBと小さく既定GOGC=100ではGC頻度が高い。
# GOGC=400でGC回数を削減(best 90,373→96,081)。swap無しのためGOMEMLIMITで上限を設けOOMを防ぐ。
# 注: GOGC=off+GOMEMLIMITはRSSが512MiB上限まで膨らみGCスラッシングで悪化したため不採用。
Environment=GOGC=400
Environment=GOMEMLIMIT=512MiB
